# nextgenans/roles/all-bootstrap/files/rc.local
#!/bin/sh -e
# This file is managed by Ansible.
# Any manual changes may be overwritten at any time.

# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
exit 0